# profile.d/40-psreadline.ps1
if (Get-Module -ListAvailable PSReadLine) {
  try {
    Set-PSReadLineOption -PredictionSource History -ErrorAction Stop
    Set-PSReadLineOption -PredictionViewStyle InlineView
  }
  catch {
    # Some non-interactive hosts do not support VT/predictions; ignore there.
  }

  Set-PSReadLineKeyHandler -Key Tab -Function MenuComplete
  Set-PSReadLineKeyHandler -Key UpArrow -Function HistorySearchBackward
  Set-PSReadLineKeyHandler -Key DownArrow -Function HistorySearchForward
}