# profile.d/50-dotnet.ps1
function db {
  dotnet build
}

function dt {
  dotnet test
}

function dr {
  dotnet run
}