# profile.d/51-navigation.ps1
function b {
  Set-Location ..
}

function bb {
  Set-Location ../..
}

function h {
  Set-Location $HOME
}

function ps-config {
  Set-Location $HOME/Documents/PowerShell
}

# TEMP - WAO  proj
function wao {
  Set-Location "$HOME\University\wao-project"
}
function wao-artifacts {
  Set-Location "$HOME\University\wao-project\artifacts"
}
function wao-gate {
  Set-Location "$HOME\University\wao-project\cinema-api-gateway"
}
function wao-audit {
  Set-Location "$HOME\University\wao-project\cinema-audit-service"
}
function wao-auth {
  Set-Location "$HOME\University\wao-project\cinema-auth-service"
}
function wao-booking {
  Set-Location "$HOME\University\wao-project\cinema-booking-service"
}
function wao-k8s {
  Set-Location "$HOME\University\wao-project\cinema-k8s"
}
function wao-movie {
  Set-Location "$HOME\University\wao-project\cinema-movie-service"
}
function wao-notification {
  Set-Location "$HOME\University\wao-project\cinema-notification-service"
}
function wao-ui {
  Set-Location "$HOME\University\wao-project\cinema-ui"
}
function wao-report {
  Set-Location "$HOME\University\wao-project\report"
}

function shower-thoughts {
  Set-Location "$HOME\.shower-thoughts"
}