# profile.d/51-navigation.ps1
function b {
  Set-Location ..
}

function bb {
  Set-Location ../..
}

function h {
  Set-Location $HOME
}

function wao {
  Set-Location "$HOME\University\wao-project"
}

function shower-thoughts {
  Set-Location "$HOME\.shower-thoughts"
}