# profile.d/52-tools.ps1
function c {
  code .
}

function c-ps {
  code $HOME/Documents/PowerShell
}

function lg {
  lazygit
}

function cl {
  Clear-Host
}

function folder {
  explorer.exe .
}

function glowE {
  glow "$HOME\.extension-list.md"
}